# benchmarks/goblin-core-replica-test.service
[Unit]
Description=Goblin Core external-logging replica test
Requires=redpanda.service
After=network-online.target redpanda.service
Wants=network-online.target

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/mnt/local/goblin-core-recovery-test
ExecStart=/usr/local/bin/goblin-core --port 6380 --numa eno1 --load /mnt/local/goblin-core/state/goblin-replica.snapshot --kafka kafka://192.168.1.49:9092/goblin-core-replication --replica-tcp 127.0.0.1 6379
Restart=on-failure
RestartSec=1s
TimeoutStopSec=30s
LimitMEMLOCK=infinity
UMask=0027

[Install]
WantedBy=multi-user.target
